# deployment/launchd/com.homelab.spicy-claude.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.homelab.spicy-claude</string>

    <key>Program</key>
    <string>/opt/homebrew/bin/node</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/edwardhallam/projects/spicy-claude/backend/dist/cli/node.js</string>
        <string>--port</string>
        <string>3002</string>
        <string>--host</string>
        <string>0.0.0.0</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/edwardhallam/projects/spicy-claude/backend</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/edwardhallam/Library/Logs/spicy-claude/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/edwardhallam/Library/Logs/spicy-claude/stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
