# deployment/launchd/com.homelab.spicy-claude-updater.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.homelab.spicy-claude-updater</string>

    <key>Program</key>
    <string>/Users/edwardhallam/projects/spicy-claude/deployment/auto-updater.sh</string>

    <key>WorkingDirectory</key>
    <string>/Users/edwardhallam/projects/spicy-claude</string>

    <key>StartInterval</key>
    <integer>600</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/edwardhallam/Library/Logs/spicy-claude-updater/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/edwardhallam/Library/Logs/spicy-claude-updater/stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>ThrottleInterval</key>
    <integer>300</integer>
</dict>
</plist>
